# packaging/systemd/powergateway-web.service
[Unit]
Description=PowerGateway Weboberflaeche
After=network.target powergateway.service
Wants=powergateway.service

[Service]
Type=simple
User=powergateway
Group=powergateway
ExecStart=/bin/sh -c 'if [ -x /opt/powergateway/venv/bin/gunicorn ]; then exec /opt/powergateway/venv/bin/gunicorn --workers 2 --threads 2 --bind 0.0.0.0:8080 --access-logfile - --error-logfile - webapp:app; else exec /usr/bin/gunicorn --workers 2 --threads 2 --bind 0.0.0.0:8080 --access-logfile - --error-logfile - webapp:app; fi'
WorkingDirectory=/opt/powergateway/src
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/etc/powergateway /opt/powergateway
ReadWritePaths=/var/lib/powergateway

[Install]
WantedBy=multi-user.target
